# courseware/config/autostart/lxterminal.desktop
[Desktop Entry]
Name=LXTerminal
Exec=lxterminal
Type=Application
